# infra/systemd/gg-docker-cleanup.service
[Unit]
Description=GratisGIS docker cleanup - prune build cache and dangling images
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Cap the buildx cache at 10 GB. Preserves recent layers for fast
# incremental rebuilds, sweeps everything older. --force is needed
# because the command otherwise prompts.
ExecStart=/usr/bin/docker builder prune --keep-storage 10GB --force
# Drop dangling and unused images. -a removes every image not
# referenced by a running OR stopped container; we run -af because
# our compose stack is `restart: unless-stopped`, so legitimate
# images are always referenced. Anything else is leftover from a
# rebuild cycle.
ExecStart=/usr/bin/docker image prune -af
ExecStart=/bin/sh -c "echo gg-docker-cleanup completed at $(date -Iseconds); docker system df"
