# infra/systemd/gg-reset-demo.timer
[Unit]
Description=Daily reset timer for the GratisGIS public test instance (#138)
Documentation=https://github.com/palavido-dev/gratis-gis

[Timer]
# 04:00 UTC daily. Chosen because it's late-night in the Americas
# (midnight ET) and pre-business-hours in Europe and Africa, so the
# ~30 - 60 second downtime hits the smallest active-user window.
# Persistent=true would catch a missed run on next boot (e.g. host
# rebooted past 04:00), which is desirable: better a reset
# happens at boot time than not at all on a day testers are abusing
# the instance.
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
# Don't run while the deploy script is actively rolling images.
# infra/deploy.sh writes /var/lib/gratis-gis-golden/deploy.lock
# at start and removes it at end; the unit refuses to run if the
# lock exists. (Implementation in deploy.sh; if the file is missing,
# the unit runs as normal.)
Unit=gg-reset-demo.service

[Install]
WantedBy=timers.target
